# packt.timer
[Unit]
Description=Download Packt free learning ebooks daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
